# remote-metal/_opt/systemd/system/run-deploy-socket.path
[Unit]
Description=Run-deploy socket service path activator

[Path]
PathModified=/tmp/run-deploy.path