# remote-incus/_opt/systemd/system/run-deploy-socket.path
[Unit]
Description=Run-deploy socket service path activator

[Path]
PathModified=/tmp/run-deploy.path

[Install]
WantedBy=multi-user.target