# remote-incus/_opt/systemd/system/run-deploy-socket.service
[Unit]
Description=Run-deploy socket service

[Service]
Type=oneshot
User=root
ExecStart=/opt/run-deploy/bin/run-deploy-socket recv